# supervisor-buildbot-master.conf
[program:buildbot]
command=/usr/bin/buildbot --verbose start --nodaemon /opt/buildbot/
autorestart=true
