# supervisor-buildbot-slave.conf
[program:buildbot-slave]
command=/usr/bin/buildslave --verbose start --nodaemon /opt/buildbot-slave/
autorestart=true
